# provision/starphleet_nuc_bootstrap.conf

description "Starphleet NUC bootstrapper, installs local SSD"

start on filesystem and net-device-up IFACE!=lo

task
script
  if mount | grep /dev/sda; then
    #this case we are already on the onboard SSD
    btrfs filesystem resize max /
    /bootstrap
  elif mount | grep /dev/sdb; then
    #this case we are coming in from the bootable USB
    /bootstrap-nuc
  fi
end script
