# overlay/etc/init/new_headquarters_orders.conf

description "Starphleet Headquarters Orders"

start on new_headquarters_orders

script

  echo "I am so ordered"

end script
