# overlay/etc/init/starphleet_dashboard.conf
description "Starphleet Ship monitoring dashboard"

start on started starphleet and started starphleet_logstream
stop on stopping starphleet

respawn

#so, this is not in a container, needs access to the host machine for stats

pre-start script
  source `which tools`
  npm ${NPM_FLAGS} install -g "${DASHBOARD_REMOTE}"
end script

script
  source `which tools`
  export PORT=3030
  exec starphleet-dashboard
end script

