# overlay/etc/init/starphleet_logstream.conf
description "Starphleet Ship log stream"

start on started starphleet
stop on stopping starphleet

respawn

#so, this is not in a container, needs access to the host machine

pre-start script
  cd /var/starphleet
  export HOME=$(pwd)
  npm install "git+https://github.com/wballard/tailor.git"
end script

script
  cd /var/starphleet/node_modules/tailor
  export HOME=$(pwd)
  ./bin/tailor /var/log/upstart 3031 > /dev/null
end script

post-stop script
  set +e
  umount -f /var/log/upstart
end script
