# overlay/etc/init/starphleet_dashboard.conf
description "Starphleet Ship monitoring dashboard"

start on started starphleet and started starphleet_logstream
stop on stopping starphleet

respawn

#so, this is not in a container, needs access to the host machine for stats

pre-start script
  source `which tools`
  [ -d "${DASHBOARD_LOCAL}" ] || mkdir -p "${DASHBOARD_LOCAL}"
  export HOME="${ADMIRAL_HOME}"
  starphleet-git-synch "${DASHBOARD_REMOTE}" "${DASHBOARD_LOCAL}" || true
  cd "${DASHBOARD_LOCAL}"
  #fuck NPM not getting uglify-js 2.4.5 timely
  npm install uglify-js@2.4.4
  npm install
end script

script
  source `which tools`
  ulimit -n ${MAX_OPEN_FILES}
  export HOME="${ADMIRAL_HOME}"
  export PORT=3030
  cd "${DASHBOARD_LOCAL}"
  exec npm start
end script

