# overlay/etc/init/starphleet_dashboard.conf
description "Starphleet Ship monitoring dashboard"

start on started starphleet
stop on stopping starphleet

respawn

#so, this is not in a container, needs access to the host machine for stats

pre-start script
  source `which tools`
  cd /var/starphleet
  echo 'building' > "/var/starphleet/dashboard.starphleetstatus"
  export HOME=$(pwd)
  npm install "${DASHBOARD_NPM}"
  echo 'built' > "/var/starphleet/dashboard.starphleetstatus"
end script

script
  source `which tools`
  cd /var/starphleet/node_modules/starphleet.dashboard
  export HOME=$(pwd)
  npm start
end script

post-start script
  echo 'online' > "/var/starphleet/dashboard.starphleetstatus"
end script

post-stop script
  echo 'offline' > "/var/starphleet/dashboard.starphleetstatus"
end script
