# overlay/etc/init/starphleet_rebuild_dead_containers.conf
description "Starphleet Rebuild Dead Containers"

start on started starphleet
stop on stopping starphleet

respawn

script
  /etc/init/starphleet_rebuild_dead_containers.start
end script
