# overlay/etc/cron.daily/restart_nginx
#!/bin/sh

/usr/bin/starphleet-restart-nginx
